# misc/etc/xdg/autostart/im-launch.desktop
[Desktop Entry]
Type=Application
Name=IM Launch
Exec=im-launch true
NoDisplay=true
X-Deepin-MaximumRAM=200
